# crontab
0 1 * * * python /usr/src/app/gnucash-updater.py
